# host/scripts/im.deepblue.cloudflare-ddns.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>im.deepblue.cloudflare-ddns</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/deepblue/services/scripts/cloudflare_ddns.py</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/deepblue/services/nginx/logs/ddns_error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/deepblue/services/nginx/logs/ddns_out.log</string>
</dict>
</plist>
